# etc/systemd/system/trackdirect_backend.service
[Unit]
Description=Trackdirect data collector backend

[Service]
ExecStart=
ExecReload=
Type=notify
Restart=always
